# com.bhoffman.cardiology-tracker-push.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.bhoffman.cardiology-tracker-push</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/benjaminhoffman/Documents/Claude/Research/Award Funding Tracker/cardiology-tracker/auto-push.sh</string>
    </array>

    <!-- Run every Monday at 8:30 AM, after the Cowork task at 8:09 AM -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/cardiology-tracker-push.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/cardiology-tracker-push.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
